# debian/postrm
#!/bin/sh
set -e
case "$1" in
    remove)
    rm -f /usr/lib/python2.7/dist-packages/picam
    update-rc.d -f picam remove
    ;;

esac
exit 0
