# debian/postrm
#!/bin/sh
set -e
case "$1" in
    remove)
    update-rc.d -f picam remove
    ;;

esac
exit 0
